# minecraftreboot.timer
[Unit]
Description=MineCraft Server Reboot

[Timer]
OnCalendar=Mon *-*-* 5:00:00
Persistent=false

[Install]
WantedBy=timers.target

